# dashboard/app_5174/apptainer/app/gedit/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:vncserver]
command=/start-gedit.sh
autostart=true
autorestart=true
stderr_logfile=/tmp/vncserver.err.log
stdout_logfile=/tmp/vncserver.out.log
priority=10

[program:gedit]
command=sh -c "sleep 5 && DISPLAY=:1 gedit"
autostart=true
autorestart=true
stderr_logfile=/tmp/gedit.err.log
stdout_logfile=/tmp/gedit.out.log
priority=20
